# bin/agent/mss-agent.service
[Unit]
Description=Mandriva Server Setup Agent
After=syslog.target network.target nss-lookup.target

[Service]
ExecStart=/usr/sbin/mss-agent

[Install]
WantedBy=multi-user.target
