# bin/agent/mss-agent.service
[Unit]
Description=Mandriva Server Setup Agent

[Service]
ExecStart=/usr/sbin/mss-agent

[Install]
WantedBy=multi-user.target
